# transfermeit.socket
[Socket]
ListenStream = 8082

[Install]
WantedBy=sockets.target